# monitor.service
[Unit]
Description=Simple Server Monitor
After=network.target nginx.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/monitor/monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
